# vision360/scrapper-api/deploy/planning-sync.service
[Unit]
Description=Planning sync (Aurion -> DB)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/scrapperplanning.env
ExecStart=/bin/bash /opt/scrapperPlanning/scripts/planning-sync.sh
